# nextcloud-agent-system/config/nextcloud-refresh.service
[Unit]
Description=Nextcloud Data Refresh Agent
After=network.target apache2.service mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nextcloud-setup/nextcloud-data-distribution
ExecStart=/opt/nextcloud-setup/nextcloud-data-distribution/.venv/bin/python3 \
    agents/agent_refresh.py \
    --env config/.env \
    --daemon \
    --interval 60
Restart=always
RestartSec=10
StandardOutput=append:/opt/nextcloud-setup/nextcloud-data-distribution/logs/refresh.log
StandardError=append:/opt/nextcloud-setup/nextcloud-data-distribution/logs/refresh_error.log

[Install]
WantedBy=multi-user.target
